# systemd_setup/vault.service
[Unit]
Description=The Vault — Post-Quantum Key Manager
After=network-online.target dev-ttyACM0.device dev-ttyAMA0.device
Wants=network-online.target dev-ttyACM0.device dev-ttyAMA0.device
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=vault
Group=vault
SupplementaryGroups=dialout
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/python3 /opt/cobratail/bin/cript_keeper.py serve
WorkingDirectory=/home/vault
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
